# src/python_exporter/python_exporter.conf
[program:python_exporter]
command=python3 exporter.py
priority=992
directory=/usr/bin/
autostart=true
autorestart=true
redirect_stderr=true
